# fanexid/rootfs/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Nginx..."
mkdir -p /run/nginx
exec nginx -g "daemon off;"
